# rootfs/etc/services.d/redis/finish
#!/usr/bin/execlineb -S0

# Cleanup logging on shutdown
rm -f /var/log/redis/*;