# rootfs/etc/services.d/redis/run
#!/usr/bin/with-contenv bash
set -e

REDIS_PASSWORD=${REDIS_PASSWORD:-}

exec s6-setuidgid redis redis-server /etc/redis.conf ${REDIS_PASSWORD:+--requirepass $REDIS_PASSWORD};